# lib/ssh_gen_keys.conf
# ssh_gen_keys - task to generate the sshd keys if they don't exist.
#
# When a VM is published as a template, the sshd keys are deleted 
# to make sure that each instance will have a unique set of sshd keys.
# This task generates those keys.
#

description     "Generates the sshd keys if they are missing"

start on filesystem or runlevel [2345]

task

script

if [ -f /etc/ssh/sshd_config -a ! -f "/etc/ssh/ssh_host_rsa_key" ]; then
  # missing host keys for the sshd; let's generate them.
  dpkg-reconfigure openssh-server
fi

end script

